# jarvis_agent/deploy/jarvis-agent.service
[Unit]
Description=Jarvis smart car patrol agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jetson
WorkingDirectory=/home/jetson/jarvis-agent
EnvironmentFile=/home/jetson/jarvis-agent/.env
ExecStart=/home/jetson/jarvis-agent/.venv/bin/python -m uvicorn jarvis_agent.api:create_app --factory --host 0.0.0.0 --port 8100
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
